# runtime/infrastructure/systemd/heiwa-nats.service
[Unit]
Description=Heiwa local NATS bus
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nats-server -js
Restart=always
RestartSec=2
User=%i
WorkingDirectory=/home/%i/heiwa-limited
StandardOutput=append:/home/%i/heiwa-limited/runtime/logs/nats.systemd.out.log
StandardError=append:/home/%i/heiwa-limited/runtime/logs/nats.systemd.err.log

[Install]
WantedBy=default.target
