# infra/local/systemd/heiwa-worker-manager.service
[Unit]
Description=Heiwa worker manager executor
After=network.target heiwa-nats.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/heiwa-limited
EnvironmentFile=/home/%i/heiwa-limited/.env.worker
Environment=HEIWA_WORKSPACE_ROOT=/home/%i/heiwa-limited
Environment=HEIWA_LLM_MODE=local_only
Environment=NATS_URL=nats://127.0.0.1:4222
ExecStart=/home/%i/heiwa-limited/cli/scripts/ops/worker_manager_daemon.sh
Restart=always
RestartSec=2
StandardOutput=append:/home/%i/heiwa-limited/runtime/logs/worker-manager.systemd.out.log
StandardError=append:/home/%i/heiwa-limited/runtime/logs/worker-manager.systemd.err.log

[Install]
WantedBy=default.target
